# ansible/roles/install/files/systemd.services/fpm-waker@.socket
[Socket]
ListenStream=/var/run/docker-apps/%i/vhost.fpm-waker.sock
SocketUser=www-data
SocketGroup=www-data

[Install]
WantedBy=sockets.target